# services/rig/before-hours.timer
[Unit]
Description=Run before hours

[Timer]
OnCalendar=Mon..Fri *-*-* 08:00:00
Unit=rig-stop.target

[Install]
WantedBy=timers.target
