# services/rig/rig-stop.target
[Unit]
Description=Stops the rig services
